# jeeves-web.service
[Unit]
Description=Jeeves Web Interface
After=network.target llama-server.service
Wants=llama-server.service

[Service]
Type=simple
User=drinrin
WorkingDirectory=/home/drinrin
ExecStart=/usr/bin/python3 /home/drinrin/jeeves_web.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
